# roles/scripts/files/systemd/user/battery_notify.service
[Unit]
Description=Battery Notification Service

[Service]
ExecStart=/home/shadow/.local/bin/battery_notify.sh

[Install]
WantedBy=default.target
